# AGI_COMPANY/infrastructure/systemd/minecraft.service
[Unit]
Description=Minecraft Server for AGI Agents
After=network.target

[Service]
WorkingDirectory=/opt/minecraft-server
ExecStart=/usr/bin/java -Xmx4G -jar paper-1.20.4-496.jar nogui
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
