# systemd/aos-brain.service
[Unit]
Description=AOS Brain - Autonomous Operations System
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=forking
User=root
WorkingDirectory=/root/.openclaw/workspace/AOS/brain
ExecStart=/usr/bin/tmux new-session -d -s aos-brain -c /root/.openclaw/workspace/AOS/brain 'python3 brain.py'
ExecStop=/usr/bin/tmux send-keys -t aos-brain C-c
ExecStopPost=/usr/bin/tmux kill-session -t aos-brain
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
